# HW2/my_test.service
[Unit]
Description=TestService

[Service]
Type=simple
ExecStart=/home/ruslan/hw2/my_service.sh

[Install]
WantedBy=multi-user.target
